# python_server_service.conf
description "Python Server Service"
author "Jonathan De Leon"
  
start on (local-filesystems and net-device-up)
stop on runlevel [06]
    
respawn
exec /var/www/python_server/server.py
